# adb-deamon/adb-autoconnect.service
[Unit]
Description=ADB Auto Connect and Port Forward Service
After=network-online.target multi-user.target
Wants=network-online.target

[Service]
Type=oneshot
# Script path (needs to be modified according to the actual installation location)
ExecStart=/usr/local/bin/adb-autoconnect.sh
# Retry on failure
Restart=on-failure
RestartSec=30
# Standard output and error output
StandardOutput=journal
StandardError=journal
# User to run as (uncomment and modify if running as a specific user is required)
# User=your_username
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

